# ticketwatch-watcher.timer
[Unit]
Description=TicketWatch Watcher Timer (every 5 minutes)
Requires=ticketwatch-watcher.service

[Timer]
# Run every 5 minutes
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
